# etc/local.d/mysqlramfs.start
#!/bin/bash

if mountpoint -q /var/ramfs; then
	echo "copying mysql into ramfs"
	cp -R /var/lib/mysql /var/ramfs
	chown -R mysql:mysql /var/ramfs/mysql

	if [ -f /etc/mysql/my_ram.cnf ]; then
		rm /etc/mysql/my.cnf
		ln -s /etc/mysql/my_ram.cnf /etc/mysql/my.cnf
	fi
fi


